# packages/app-admin/puppet/files/systemd/puppet-agent.service
[Unit]
Description=puppet agent daemon

[Service]
ExecStart=/usr/bin/puppet agent --no-daemonize
Restart=always

[Install]
WantedBy=multi-user.target
